# kube/etc/services.d/php-fpm/run
#!/command/execlineb -P
with-contenv
exec php-fpm82
